# 80-local.rules
SUBSYSTEMS=="usb", ACTION=="remove", ATTRS{serial}=="your_device_id_here", RUN+="/home/username/bin/shutdown.sh"

